# conf/supervisor/conf.d/pwmfan.conf
[program:tlsmFanControl]
command=/home/pi/tlsm/fan-ctl.sh
directory=/home/pi/tlsm/
autorestart=true
startretries=1000
stderr_logfile=/var/log/fan-ctl-err.log
stdout_logfile=/var/log/fan-ctl-out.log
